# sources/meta-sysworxx/recipes-core/di-setup/files/di-setup.service
[Unit]
Description=Setup driver for digital inputs

[Service]
Type=oneshot
ExecStart=/usr/sbin/di-setup.sh

[Install]
WantedBy=multi-user.target
