# sources/meta-sysworxx/recipes-core/util-linux/files/fstrim.service
[Unit]
Description=Discard unused blocks on filesystems

[Service]
Type=oneshot
ExecStart=/usr/sbin/fstrim --all --verbose
