# sources/meta-sysworxx/recipes-core/adc-setup/files/adc-setup.service
[Unit]
Description=Set mode of all analog inputs

[Service]
Type=oneshot
ExecStart=/usr/sbin/adc-setup.sh

[Install]
WantedBy=multi-user.target
